# soft-hwclock-tick.service
[Unit]
Description=Save hwclock by timer

[Service]
Type=oneshot
ExecStart=/opt/soft-hwclock/soft-hwclock tick > /dev/null

